# dockerfiles/supervisord.conf
[supervisord]
nodaemon=true
[program:main]
directory=/compute-backend
command=node index.js； manage 进程

