# dockerfiles/supervisord.conf
[supervisord]
nodaemon=true

[program:main]
directory=/compute-backend
command=node index.js；

[program:matdaemon]
directory=/compute-backend/service
command=python matdeamon.py prod；