# dockerfiles/supervisord.conf
[supervisord]
nodaemon=true

[program:main]
command=node /app/index.js

[program:matdaemon]
command=python /app/services/matdaemon.py